# debian-pkg/DEBIAN/postrm
#!/bin/bash

systemctl disable --now janus.service > /dev/null 2>&1 || true
